# templates/.claude/skills/git-report/resources/com.user.claude-git-report.template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- 任务标识 -->
    <key>Label</key>
    <string>com.user.claude-git-report</string>

    <!-- 任务描述 -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{PROJECT_ROOT}}/templates/.claude/skills/git-report/resources/git-report-automation.sh</string>
    </array>

    <!-- 定时执行：每周一上午8:00 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <!-- 工作目录 -->
    <key>WorkingDirectory</key>
    <string>{{PROJECT_ROOT}}</string>

    <!-- 标准输出和错误输出 -->
    <key>StandardOutPath</key>
    <string>{{HOME}}/.claude/logs/daily-report-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/.claude/logs/daily-report-launchd-error.log</string>

    <!-- 环境变量 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- 运行条件 -->
    <key>RunAtLoad</key>
    <false/>

    <!-- 保持活跃 -->
    <key>KeepAlive</key>
    <false/>

    <!-- 退出时重启 -->
    <key>AbandonProcessGroup</key>
    <true/>
</dict>
</plist>